# root/etc/uci-defaults/luci-tailscale
#!/bin/sh

# Run tailscale for the first time
/etc/init.d/tailscale start
/etc/init.d/tailscale up
/etc/init.d/tailscale down
/etc/init.d/tailscale ip
/etc/init.d/tailscale status
/etc/init.d/tailscale web

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@tailscale[-1]
	add ucitrack tailscale
	set ucitrack.@tailscale[-1].init=tailscale
	commit ucitrack
EOF


# 重载 rpcd 服务，以便新安装的界面能顺利显示，否则就必须关闭浏览器重新登录才能显示新安装的界面菜单
/etc/init.d/rpcd reload
exit 0